# scripts/com.itrimble.dotfiles.maintenance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.itrimble.dotfiles.maintenance</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ian/.dotfiles/scripts/weekly-maintenance.sh</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>StandardOutPath</key>
    <string>/Users/ian/.dotfiles/logs/launchd-stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/ian/.dotfiles/logs/launchd-stderr.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/ian/.dotfiles</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/ian/.local/bin:/Users/ian/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/local/sbin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/ian</string>
        <key>USER</key>
        <string>ian</string>
    </dict>
</dict>
</plist>
